# machines/eb-reverse-proxy/etc/systemd/system/certbot.service.d/override.conf
[Service]
ExecStartPost=/usr/local/sbin/certbot-post.sh
